# deploy/haive.service
[Unit]
Description=HAIVE Nuxt server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/haive
EnvironmentFile=/home/pi/haive/.env
ExecStart=/usr/bin/node /home/pi/haive/.output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1

[Install]
WantedBy=multi-user.target
